# systemd/suedtirolmobil-bot.service
[Unit]
Description=suedtirolmobilAI Telegram Bot
After=network.target suedtirolmobil.service

[Service]
Type=simple
# adjust this path to the repository location
WorkingDirectory=/path/to/suedtirolmobilAI
ExecStart=/path/to/suedtirolmobilAI/venv/bin/python -m src.telegram_bot --api-url http://localhost:8000
EnvironmentFile=/path/to/suedtirolmobilAI/.env
Restart=on-failure

[Install]
WantedBy=multi-user.target
